# server/supervisor-video-detection.conf
[program:geprekin-video-detection]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/report.geprekincloud.tech/artisan queue:work redis --queue=video-detection --sleep=3 --tries=1 --timeout=1800
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/report.geprekincloud.tech/storage/logs/video-detection-worker.log
stopwaitsecs=1900
